# systemd/mail-notify.timer
[Unit]
Description=Mail notifier or libnotify

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
